# deploy/mon.conf
start on started network-services
stop on stopping network-services
respawn
chdir /var/www/mon/
exec /var/www/mon/main.linux